# usermanager.service
[Unit]
Description=UserManager — Centralised user identity and auth service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sedjwin
Group=sedjwin
WorkingDirectory=/home/sedjwin/UserManager
ExecStart=/home/sedjwin/UserManager/.venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8005 \
    --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=usermanager

[Install]
WantedBy=multi-user.target
